# systemd/telegram-zoomer.service
[Unit]
Description=Telegram NYT to Russian Zoomer Translator Bot
After=network.target

[Service]
Type=simple
User=yourusername
WorkingDirectory=/path/to/telegram_zoomer
ExecStart=/usr/bin/python3 /path/to/telegram_zoomer/bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-zoomer
Environment=PYTHONUNBUFFERED=1
# Load environment variables from file
EnvironmentFile=/path/to/telegram_zoomer/.env

[Install]
WantedBy=multi-user.target 